# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true

[program:jenkins]
command=./run.sh
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:build-groovy-jobs.sh]
command=./build-groovy-jobs.sh
startsecs = 0
autorestart = false
exitcodes=0
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:execute-jenkins-cli-commands.sh]
command=./execute-jenkins-cli-commands.sh
startsecs = 0
autorestart = false
exitcodes=0
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[group:all]
programs=nginx,jenkins,build-groovy-jobs.sh

